# rootfs/build/services/vmd.conf
[program:vmd]
; The single service: vmd owns QEMU + swtpm + the web (noVNC / serial console / power API), driven
; by /vms/vmd.toml. Toggled with ubase ZSRV_vmd (see supervisor/).
command=/build/bin/vmd run
autostart=true
; exit 0 = clean guest poweroff (stay off); non-zero = crash (restart) — vmd propagates QEMU's code
autorestart=unexpected
exitcodes=0
startsecs=10
stopsignal=TERM
; vmd performs the ACPI shutdown + grace itself; keep below the container stop_grace_period
stopwaitsecs=160
redirect_stderr=true
stdout_logfile=/var/log/supervisord-vmd.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
